# systemd/openclaw-concert-marketing.service
[Unit]
Description=OpenClaw Concert Marketing Worker (Consolidated Fleet)
Documentation=https://docs.openclaw.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=concert-marketing
Group=concert-marketing
WorkingDirectory=/home/concert-marketing
EnvironmentFile=-/etc/openclaw/concert-marketing.env
ExecStart=/usr/bin/openclaw gateway run
Restart=always
RestartSec=8s
TimeoutStopSec=20s

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/concert-marketing/.openclaw /srv/consolidated-fleet /tmp

[Install]
WantedBy=multi-user.target
